# nexus433.service.in
# Copy this to /etc/systemd/system/

[Unit]
Description=433MHz sensor receiver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=@NEXUS433_BINARY_PATH@

[Install]
WantedBy=multi-user.target
